# pmsb.conf
[program:pmsb-web]
command=/home/pmsb/bin/gunicorn_start
user=pmsb
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/pmsb/logs/gunicorn-error.log